# jetson-nano-2gb-devkit/meta-docker-app/recipes-containers/docker-app/files/start-stop-webapp.service
[Unit]
Description=Start Stop WebApp Service

[Service]
Type=simple
ExecStart=/etc/init.d/start-stop-webapp

[Install]
WantedBy=multi-user.target

